# build_vagrant/rq_supervisor.conf
[program:myworker]
; Point the command to the specific rq command you want to run.
; If you use virtualenv, be sure to point it to
; /path/to/virtualenv/bin/rq
; Also, you probably want to include a settings module to configure this
; worker.  For more info on that, see http://python-rq.org/docs/workers/
command=/vagrant/venv/bin/rq worker -c rq_config high normal low
process_name=%(program_name)s_%(process_num)02d

; If you want to run more than one worker instance, increase this
numprocs=4

; This is the directory from which RQ is ran. Be sure to point this to the
; directory where your source code is importable from
directory=/vagrant/

; RQ requires the TERM signal to perform a warm shutdown. If RQ does not die
; within 10 seconds, supervisor will forcefully kill it
stopsignal=TERM

; These are up to you
autostart=true
autorestart=true

; Logs
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log